# scripts/crontab
*/1 * * * * /bin/sh -lc /scripts/data_sync.sh