# scripts/crontab
# Sync data to R2 every minute
*/10 * * * * /scripts/data_sync.sh
