# scripts/crontab
# Sync data to R2 every minute
*/1 * * * * /scripts/data_sync.sh
